# dot_config/autostart/ckb-next.desktop
[Desktop Entry]
Type=Application
Name=ckb-next
Exec=ckb-next --background
