# dot_config/autostart/firefox.desktop
[Desktop Entry]
Type=Application
Name=Firefox default
Exec=firefox -P default
