# dot_config/autostart/ghostty.desktop
[Desktop Entry]
Type=Application
Name=ghostty
Exec=ghostty
