# dot_config/autostart/remmina-applet.desktop
[Desktop Entry]
Type=Application
Name=Remmina Applet
Exec=remmina -i
