# dot_config/autostart/slack.desktop
[Desktop Entry]
Type=Application
Name=Slack
Exec=slack
