# dot_config/autostart/udiskie.desktop
[Desktop Entry]
Type=Application
Name=udiskie
Exec=udiskie
